# etc/udev/rules.d/65-storage-optimization.rules
# ============================================================================
# Storage Device Optimization - Speed and reliability
# ============================================================================
# Optimizes SSD/HDD settings for better responsiveness

# NVMe: Enable Autonomous Power State Transition (APST)
# Allows drive to automatically enter deeper sleep states
# Reduces power consumption on idle with no performance impact
ACTION=="add", SUBSYSTEM=="block", KERNEL=="nvme*", \
  RUN+="/bin/bash -c 'nvme set-feature -n /dev/$name -f 2 -v 2 2>/dev/null || true'"

# NVMe: Set LBA format (if available) for better performance
# Comment out if causing boot delays
# ACTION=="add", SUBSYSTEM=="block", KERNEL=="nvme*", \
#   RUN+="/bin/bash -c 'nvme format /dev/$name --lbaf=0 2>/dev/null || true'"

# SSD: Increase read-ahead for better sequential performance
# Helps with: fast file transfers, media playback, backups
# 256KB is good balance between performance and memory usage
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z]*|nvme*", ATTR{queue/rotational}=="0", \
  ATTR{queue/read_ahead_kb}="256"

# HDD: Conservative read-ahead to avoid excessive seeking
# HDDs have high seek time, so aggressive read-ahead can hurt random access
# 128KB is safe for mechanical drives
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z]*", ATTR{queue/rotational}=="1", \
  ATTR{queue/read_ahead_kb}="128"

# Enable TRIM support for SSDs (for garbage collection)
# Helps maintain SSD performance over time
# Works with ext4, btrfs, f2fs, ntfs
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z]*", \
  RUN+="/usr/bin/hdparm -M 128 /dev/$name 2>/dev/null || true"

# Set NCQ/AHCI for SATA drives (Native Command Queuing)
# Allows drive to optimize command order for better throughput
# ~5-15% performance improvement for SATA drives
ACTION=="add", SUBSYSTEM=="ata_port", RUN+="/bin/bash -c 'echo 32 > /sys/module/libata/parameters/force_nrq 2>/dev/null || true'"

# Virtual/Cloud Block Devices: No optimization needed
# ACTION=="add", SUBSYSTEM=="block", KERNEL=="vd[a-z]*", \
#   ATTR{queue/read_ahead_kb}="512"
