# etc/udev/rules.d/66-input-optimization.rules
# ============================================================================
# Input Device Optimization - Perfect for keyboard-driven WMs
# ============================================================================
# Optimizes keyboards, trackpads, and input for chadwm/ohmychadwm

# Increase keyboard repeat rate (useful for vim navigation)
ACTION=="add", SUBSYSTEM=="input", KERNEL=="event*", ATTRS{name}=="*Keyboard*", \
  RUN+="/bin/bash -c 'echo 200 > /proc/sys/kernel/autorepeat_interval 2>/dev/null || true'"

# Trackpad: disable while typing (prevent accidental clicks)
ACTION=="add", SUBSYSTEM=="input", ATTRS{name}=="*Touchpad*|*TrackPoint*", \
  ENV{LIBINPUT_IGNORE_DUPLICATE_EVENT}="1"

# Mouse: higher polling rate for responsive pointer
ACTION=="add", SUBSYSTEM=="input", KERNEL=="mice", \
  RUN+="/bin/bash -c 'echo 1000 > /sys/module/usbhid/parameters/jitter 2>/dev/null || true'"

# Reduce mouse acceleration for precise clicking (important for tiling WM)
ACTION=="add", SUBSYSTEM=="input", ATTRS{name}=="*Mouse*", \
  ENV{MOUSE_ACCEL}="0", ENV{MOUSE_SPEED}="1.0"

