# etc/udev/rules.d/68-sound-power.rules
# ============================================================================
# Advanced Sound Card Power Management
# ============================================================================
# Fine-tuned for audio quality + power efficiency

# ALSA: prevent clicks/pops from power transitions
ACTION=="add", SUBSYSTEM=="sound", DRIVER=="snd_hda_intel", \
  ATTR{model}="auto", ATTR{power_save}="10"

# Jack detection for headphones (auto-disable speakers)
ACTION=="add", SUBSYSTEM=="sound", KERNEL=="card*", \
  RUN+="/usr/bin/alsactl restore 2>/dev/null || true"
